# api/supervisord.conf
; ============================================================================
; supervisord.conf — WWA API (VPS Hostinger)
;
; Installation :
;   apt install supervisor
;   cp api/supervisord.conf /etc/supervisor/conf.d/wwa.conf
;   # Remplacer /var/www/wwa/api par le chemin réel
;   supervisorctl reread && supervisorctl update
;   supervisorctl start wwa:*
;
; Scalabilité Phase 2 (Redis + multi-instances) :
;   - Passer QUEUE_CONNECTION=redis dans .env
;   - Passer REVERB_SCALING_ENABLED=true + REDIS_HOST dans .env
;   - numprocs=4 sur wwa-queue si volume élevé
; ============================================================================

[group:wwa]
programs=wwa-reverb,wwa-queue

; ── Reverb WebSocket Server ──────────────────────────────────────────────────
; Long-running : écoute les connexions WS et dispatche les broadcasts.
[program:wwa-reverb]
command=php /var/www/wwa/api/artisan reverb:start --host="0.0.0.0" --port=8080 --no-interaction
directory=/var/www/wwa/api
user=www-data
autostart=true
autorestart=true
startretries=5
numprocs=1
; Redémarrage propre : tue les processus enfants
stopasgroup=true
killasgroup=true
; Délai avant SIGKILL après SIGTERM
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/wwa-reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_redirect=true
; Évite les logs dupliqués
redirect_stderr=true

; ── Queue Worker ─────────────────────────────────────────────────────────────
; Traite les jobs async : ProcessAttachment (validation MIME), emails, etc.
; Les notifications de message (NewMessageNotification) sont synchrones —
; ce worker gère les autres jobs.
[program:wwa-queue]
command=php /var/www/wwa/api/artisan queue:work --sleep=3 --tries=3 --timeout=30 --max-time=3600 --no-interaction
directory=/var/www/wwa/api
user=www-data
autostart=true
autorestart=true
startretries=5
numprocs=1
stopasgroup=true
killasgroup=true
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/wwa-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true

; ── Note nginx (reverse proxy pour Reverb) ───────────────────────────────────
; Ajouter dans le bloc server nginx du domaine API :
;
; location /app {
;     proxy_pass             http://127.0.0.1:8080;
;     proxy_http_version     1.1;
;     proxy_set_header       Upgrade $http_upgrade;
;     proxy_set_header       Connection "upgrade";
;     proxy_set_header       Host $host;
;     proxy_set_header       X-Real-IP $remote_addr;
;     proxy_cache_bypass     $http_upgrade;
;     proxy_read_timeout     300s;
; }
